# transcription-server.service
[Unit]
Description=Live Transcription Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/live-transcribe-and-create-live-responses
ExecStart=/root/live-transcribe-and-create-live-responses/server_manager.sh monitor
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
